# yt-downloader.service
[Unit]
Description=YouTube Downloader Flask Application
After=network.target
Documentation=https://github.com/yourusername/yt-downloader

[Service]
Type=notify
User=yt-downloader
Group=yt-downloader
WorkingDirectory=/opt/yt-downloader

# Environment variables
Environment="PATH=/opt/yt-downloader/venv/bin"
EnvironmentFile=-/etc/yt-downloader/yt-downloader.env

# Start command - using gunicorn
ExecStart=/opt/yt-downloader/venv/bin/gunicorn \
    --config /opt/yt-downloader/gunicorn_config.py \
    app:app

# Restart policy
Restart=on-failure
RestartSec=10
StartLimitInterval=60s
StartLimitBurst=3

# Security and resource limits
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/downloads /var/log/yt-downloader

# Resource limits (adjust based on your system)
LimitNOFILE=65535
LimitNPROC=512

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=yt-downloader

[Install]
WantedBy=multi-user.target
